# files/etc/systemd/system/mark2-automount.service
[Unit]
Description=Automount USB drives with udisks2

[Service]
ExecStart=/usr/bin/udisks-automount.sh

[Install]
WantedBy=default.target
